# build/automc/buildbot.cron.weekly
#!/bin/sh

cd /home/buildbot/Maildir
rm -f commits*/cur/*
